# java/sync/debian/spushnik.postrm
#!/bin/bash
set -e

invoke-rc.d jetty status
jetty_status=$?

if [[ $jetty_status -eq 0 && -f /etc/init.d/jetty ]]; then
	invoke-rc.d jetty restart	
elif [ $jetty_status -ne 3 ]; then
	echo "Jetty not correctly shutdown.\nPlease restart it manually to stop using the spushnik webapp.\n"
fi

[ -f /etc/init.d/spushnik ] && unlink /etc/init.d/spushnik

exit 0

