# ui/rpm/SOURCES/obm-services.cron.d
#
# Regular cron jobs for the obm-services package
#
0 *	* * *	root /usr/share/obm/www/auto/ldapContacts.pl
0 0	* * *	root /usr/share/obm/www/auto/cyrusGetQuotaUse.pl
