# jetty/debian/obm-jetty.postrm
#! /bin/bash

set -e

invoke-rc.d jetty stop

# Revert the modifications applied to /etc/default/jetty
cp -a /etc/default/jetty.orig /etc/default/jetty

# Revert the mdofications applied to /etc/jetty/jetty.xml
cp -a /etc/jetty/jetty.xml.orig /etc/jetty/jetty.xml

invoke-rc.d jetty start

#DEBHELPER#

exit 0
