# java/solr/debian/obm-solr.cron.d
#
# Optimize solr index every day @midnight
#
0 0	* * *	root /usr/share/obm-solr/obm_optimize_index.py
