# ui/debian/obm-filter.cron.d
#
# Regular cron jobs for the obm-filter package
#
0 3	* * 0	root /usr/bin/update-sa-rules.sh	
