# ui/debian/obm-postfix.postrm
#!/bin/sh

#set -x

case "$1" in
	purge|remove|failed-upgrade|abort-install|abort-upgrade|disappear)

		# Supression de la conf POSTFIX
		invoke-rc.d postfix stop || true
		mv /etc/postfix/main.cf.orig /etc/postfix/main.cf || true
	;;
	

    upgrade)
    ;;

esac

#DEBHELPER#


exit 0
